# com.emailsorter.listen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.emailsorter.listen</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vladimirhohlov/email-sorter/venv/bin/python</string>
        <string>-u</string>
        <string>/Users/vladimirhohlov/email-sorter/sorter.py</string>
        <string>--listen</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vladimirhohlov/email-sorter</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/email-sorter.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/email-sorter.err</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
